# s6/sv/mount-efivars/up
#!/bin/execlineb -P
if { s6-echo "Mounting efivars..." }
foreground {
   if { test -d /sys/firmware/efi }
       if -n { mountpoint -q /sys/firmware/efi/efivars }
          mount -n -t efivarfs -o ro efivarfs /sys/firmware/efi/efivars
}
