# s6/sv/cgroups/up
#!/bin/execlineb -P
if { s6-echo "Mounting cgroups..." }
envfile /etc/s6/s6.conf
importas -iu CGROUP_MODE CGROUP_MODE
importas -iu CGROUP_CONTROLLERS CGROUP_CONTROLLERS
importas -iu HAVE_CONTROLLER1_GROUPS HAVE_CONTROLLER1_GROUPS
exec sh /etc/s6/scripts/mount-cgroups $CGROUP_MODE $CGROUP_CONTROLLERS $HAVE_CONTROLLER1_GROUPS
