# s6/base/run-image/service/s6-linux-init-early-getty/run
#!/bin/execlineb -P

/sbin/agetty -L -8 tty1 115200
