# s6/sv/mount-sysfs/up
#!/bin/execlineb -P
if { s6-echo "[   Virtual-FS   ] 2/6 : Mounting /sys" }
s6-mount -t sysfs -o nosuid,noexec,nodev sys /sys
