# s6/sv/cleantmp/up
#!/bin/execlineb -P
if { s6-echo "[ ReadWrite Mode ] 0/6 : Cleaning /tmp" }
/bin/sh -c "/etc/s6/scripts/clean_tmp.sh"
