# s6/current/run-image/service/s6-linux-init-runleveld/run
#!/usr/bin/execlineb -P

fdmove -c 2 1
fdmove 1 3
s6-ipcserver -1 -a 0700 -c 1 -- s
s6-sudod -dt30000 --
"/etc/s6/current"/scripts/runlevel
