# s6/sv/agetty5/run
#!/bin/execlineb -P
if { s6-echo "[    Services    ]     : Enabling tty5" }
exec /sbin/agetty -L -8 tty5 115200 
