# s6/sv/agetty5/run
#!/bin/execlineb -P
exec /sbin/agetty -L -8 tty5 115200
