# s6/sv/agetty5/run
#!/bin/execlineb -P
if { s6-echo -- "[  Enable VT's   ] 4/5 : Enabling tty5" }
exec /sbin/agetty -L -8 tty5 115200
