# s6/sv/mount-net/up
#!/usr/bin/execlineb -P
exec mount -a -O _netdev
