# s6/sv/udevd/run
#!/bin/execlineb -P
# Initiate udev
fdmove -c 2 1
if { s6-echo "Intializing udev daemon..."  }
exec -c
udevd
