# s6/sv/tmpfiles-dev/up
#!/bin/execlineb -P
exec tmpfiles --prefix=/dev --create --boot
