# s6/sv/tmpfiles-dev/up
#!/bin/execlineb -P
if { s6-echo "[ ReadWrite Mode ] 6/6 : Setting up tmpfiles" }
tmpfiles --prefix=/dev --create --boot
