# s6/sv/tmpfiles-dev/up
#!/bin/execlineb -P
if { s6-echo -- "[    Tmpfiles    ] Creating TMP files"  }

exec tmpfiles --prefix=/dev --create --boot
