# s6/sv/tmpfiles-dev/up
#!/bin/execlineb -P
if { s6-echo "Setting up tmpfiles..." }
tmpfiles --prefix=/dev --create --boot
