# s6/sv/agetty3/run
#!/bin/execlineb -P
if { s6-echo "Enabling tty3..." }
exec /sbin/agetty -L -8 tty3 115200 
