# s6/sv/cleanboot/up
#! /bin/execlineb -P

fdmove -c 2 1

if { s6-echo -- "[    cleanboot   ] 0/5 : Creating directories ... " }
if -t {
    if { s6-echo -- "[    cleanboot   ] 4/5 : Creating ICE-unix" }
    foreground { install -dm1777 /tmp/.X11-unix /tmp/.ICE-unix }
}
if { s6-echo -- "[    cleanboot   ] 5/5 : Removing filesytem watermarks" }
rm -fv /etc/nologin /forcecheck /forcequotacheck /fastboot
