# s6/sv/mnt-devpts/up
#!/bin/execlineb -P
if { s6-echo "[   Virtual-FS   ] 3/6 : Mounting /dev/pts" }
s6-mount -t devpts -o mode=0620,gid=5,nosuid,noexec devpts /dev/pts
