# s6/sv/modules/up
#!/bin/execlineb -P
if { s6-echo "[ Read-Only Mode ] 6/9 : Loading any kernel modules" }
sh -c "modules-load"
