# s6/sv/modules/up
#!/bin/execlineb -P
if { s6-echo "Loading any kernel modules..." }
sh -c "modules-load"
