# s6/sv/modules/up
#!/usr/bin/execlineb -P
exec modules-load
