# s6/sv/swap/up
#!/bin/execlineb -P
if { s6-echo "Turning on any swap/swapfiles..." }
swapon -a
