# s6/sv/swap/up
#!/bin/execlineb -P
if { s6-echo "[ Read-Only Mode ] 7/9 : Turning on any swap/swap-files" }
swapon -a
