# s6/sv/net-lo/up
#!/usr/bin/execlineb -P
exec ip link set up dev lo
