# s6/sv/net-lo/up
#! /bin/execlineb -P

fdmove -c 2 1

if { s6-echo -- "[     Net-lo     ] Setting up loopback device..." }
if { ip link set up dev lo }
