# s6/sv/net-lo/up
#!/bin/execlineb -P
if { s6-echo "Setting up loopback interface..." }
ip link set up dev lo
