# s6/sv/udevadm/up
#!/bin/execlineb -P
# Initiate udev devices & subsystems
fdmove -c 2 1

if { s6-echo -- "[      udev      ] 2/2 :  Setting up devices & subsystems..." }
if {
  foreground { udevadm trigger --action=add --type=subsystems }
  foreground { udevadm trigger --action=add --type=devices }
  udevadm settle
}
