# s6/sv/kermod/up
#!/bin/execlineb -P
if { s6-echo "Setting up kernel static node(s)..." }
foreground { if -n { test -d /run/tmpfiles.d } mkdir /run/tmpfiles.d }
foreground { kmod static-nodes --format=tmpfiles --output=/run/tmpfiles.d/kmod.conf }
