# s6/sv/sysctl/up
#!/bin/execlineb -P
if { s6-echo -- "[ System Control ] Setting kernel parameters " }
redirfd -w 1 /dev/null sysctl --system
