# s6/sv/sysctl/up
#!/bin/execlineb -P
if { s6-echo "[ Read-Only Mode ] 0/9 : Setting kernel parameters" }
redirfd -w 1 /dev/null sysctl --system
