# s6/sv/console/up
#!/bin/execlineb -P
if { s6-echo "Setting up the console..." }
sh -c "/etc/s6/scripts/console_set.sh"
