# s6/sv/console/up
#!/bin/execlineb -P
if { s6-echo "[ Read-Only Mode ] 2/9 : Setting up console" }
sh -c "/etc/s6/scripts/console_set.sh"
