# s6/sv/prep-dev/up
#!/bin/execlineb -P
if { s6-echo "[   Virtual-FS   ] 0/6 : Creating mount points" }
s6-mkdir /dev/shm
s6-mkdir /dev/pts
