# s6/sv/prep-dev/up
#!/bin/execlineb -P
if { s6-echo "Creating mount points..." }
foreground { 
  if -n { s6-test -d /dev/shm }
  s6-mkdir /dev/shm
}
foreground {
  if -n { s6-test -d /dev/pts }
  s6-mkdir /dev/pts
}
