# s6/current/run-image/service/s6-svscan-log/run
#!/usr/bin/execlineb -P

fdmove -c 1 2
redirfd -rnb 0 fifo
s6-log -bpd3 -- 1 t /run/uncaught-logs
