# classic/service/s6-svscan-log/run
#!/bin/execlineb -P

importas -i DEST DEST
redirfd -w 2 /dev/console
fdmove -c 1 2
redirfd -rnb 0 fifo
s6-setuidgid s6log
s6-log -bp -- t ${DEST}/uncaught-logs

