# s6/sv/mount-nonetwork/up
#! /bin/execlineb -P
fdmove -c 2 1

if { s6-echo -- "[    nonetwork   ] 1/1 : Mounting non-network filesystems" }
if { mount -a -t "nosysfs,nonfs,nonfs4,nosmbfs,nocifs" -O no_netdev }
