# s6/sv/agetty2/run
#!/bin/execlineb -P
exec /sbin/agetty -L -8 tty2 115200
