# s6/sv/mount-filesystems/up
#!/usr/bin/execlineb -P
foreground { mount -a -O no_netdev }
foreground { if { test -x /usr/bin/zfs } zfs mount -a }
