# s6/sv/console-setup/up
#!/bin/execlineb -P
if { s6-echo -- "[  Consle-setup  ] Setting VT font and keymap" }
exec sh /etc/s6/sv/console-setup/shell_up
