# s6/sv/console-setup/up
#!/bin/execlineb -P
exec sh /etc/s6/sv/console-setup/shell_up
