# s6/sv/hwclock/up
#!/bin/execlineb -P
if { s6-echo "[ Read-Only Mode ] 4/9 : Setting system clock from rtc0" }
envfile /etc/s6/s6.conf
importas -iu HARDWARECLOCK HARDWARECLOCK
foreground { if { s6-test $HARDWARECLOCK = UTC } hwclock --systz --utc --noadjfile }
foreground { if { s6-test $HARDWARECLOCK = localtime } hwclock --systz --localtime --noadjfile }
