# s6/sv/checkfs/up
envfile /etc/s6/s6.conf
importas -iu FORCECHECK FORCECHECK
ifelse -X { s6-test $FORCECHCK = yes }
 {
  redirfd -w 1 /dev/console
  if { s6-echo -- "[  Checkrootfs   ] >>>>> Check of filesystem was asked, please wait" }
  foreground { fsck -A -T -a -f noopts=_netdev }
  s6-echo -- "[  Checkrootfs   ] >>>>> Check of filesystem was asked, please wait"
}
if -t {
    fsck -A -T -a noopts=_netdev
}
