# run-image/service/s6-linux-init-early-getty/run
#!/bin/execlineb -P

/sbin/agetty 38400 tty1
