# s6/sv/sysuser/up
#!/usr/bin/execlineb -P

## Install esysusers first
#exec sysusers
