# s6/sv/agetty4/run
#!/bin/execlineb -P
if { s6-echo -- "[  Enable VT's   ] 3/5 : Enabling tty4" }
exec /sbin/agetty -L -8 tty4 115200
