# s6/sv/agetty4/run
#!/bin/execlineb -P
if { s6-echo "Enabling tty4..." }
exec /sbin/agetty -L -8 tty4 115200 
