# s6/current/run-image/service/.s6-svscan/finish
#!/usr/bin/execlineb -P

redirfd -w 2 /dev/console
fdmove -c 1 2
foreground { s6-linux-init-echo -- "s6-svscan exited. Rebooting." }
s6-linux-init-hpr -fr
