# s6/sv/remount-root/up
#!/bin/execlineb -P

# Mount Root Filesystem Read-Write

fdmove -c 2 1

if { s6-echo -- "[   remount-rw   ] Remounting as read-write"  }
if { s6-mount -o remount,rw / / }
