# s6/sv/remount-root/up
#!/usr/bin/execlineb -P
exec mount -o remount,rw /
