# s6/sv/remount-root/up
#!/bin/execlineb -P
if { s6-echo "Remounting root filesystem as read/write..." }
s6-mount -o remount,rw / / 
