# s6/sv/remount-root/up
#!/bin/execlineb -P
if { s6-echo "[ ReadWrite Mode ] 1/6 : Remounting root filesystem as rw" }
s6-mount -o remount,rw / / 
