# s6/sv/binfmt/up
#!/usr/bin/execlineb -P
foreground { sh /etc/s6/sv/binfmt/shell_up }
