# s6/sv/tmpfiles-setup/up
#!/bin/execlineb -P
if { s6-echo "Cleaning tmpfiles..." }
tmpfiles --exclude-prefix=/dev --create --remove --boot
