# s6/sv/tmpfiles-setup/up
#!/bin/execlineb -P
exec tmpfiles --exclude-prefix=/dev --create --remove --boot
