# classic/service/tty6/run
#!/bin/execlineb -P

/sbin/agetty --delay 02 -J 38400 tty6
