# s6/sv/mount-procfs/up
#!/bin/execlineb -P
if { s6-echo "Mounting /proc..." }
s6-mount -t proc -o nosuid,noexec,nodev proc /proc
