# s6/sv/mount-procfs/up
#!/bin/execlineb -P
foreground {
    if -n { test -e /proc/mounts }
        if -n { mountpoint -q /proc } mount -t proc proc /proc -o nosuid,noexec,nodev
}
