# s6/sv/agetty6/run
#!/bin/execlineb -P
if { s6-echo -- "[  Enable VT's   ] 5/5 : Enabling tty6" }
exec /sbin/agetty -L -8 tty6 115200
