# s6/sv/agetty6/run
#!/bin/execlineb -P
if { s6-echo "[    Services    ]     : Enabling tty6" }
exec /sbin/agetty -L -8 tty6 115200 
