# s6/sv/dmesg/up
#!/bin/execlineb -P
if { s6-echo "[ ReadWrite Mode ] 2/6 : Logging kernel boot" }
pipeline { dmesg } s6-setuidgid s6log exec -c s6-log -b -- n3 s2000000 T /var/log/dmesg
